# files/auto_share.service
[Unit]
Description=Check the network mount
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/auto_share.sh
